# src/fastnetmon_init_script_gentoo
#!/sbin/runscript

FASTNETMON="/opt/fastnetmon/fastnetmon"
PIDFILE="/var/run/fastnetmon.pid"
ARGS="--daemonize"

depend() {
    need net
}

start() {
    ebegin "Starting fastnetmon"
    if [ -f "${PIDFILE}" ]; then
        einfo "  Removing stale pidfile ${PIDFILE}"
        rm -f "${PIDFILE}" 1>/dev/null
    fi

    $FASTNETMON $ARGS
    eend $?
}

stop() {
    ebegin "Stopping fastnetmon"

    if [ -f $PIDFILE ]; then
    	kill -9 $(cat $PIDFILE) 2>/dev/null
	RETVAL=$?
    fi
    if [ -n $RETVAL ] && [ "$RETVAL" -ne "0" ]; then
    	ACTUAL_PID=$(pidof fastnetmon)
	[ -n $ACTUAL_PID ] && kill -9 $ACTUAL_PID 2>/dev/null
	RETVAL=$?
    fi

    rm -f "${PIDFILE}"
    eend $RETVAL
}